# software/docker/supervisord.ini
[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx-php]
command=/usr/bin/php /var/www/start.php